# stratumd.conf
description "stratum daemon"
author "coelacanthus <coelacanth-soft@nifty.com>"

start on runlevel [2345]
stop on runlevel [016]

chdir /path/to/yiimp_install_scrypt
exec python stratumd.py


